# feishu-bot.service
# Nono CoWork - Feishu Bot systemd service
# Copy to /etc/systemd/system/feishu-bot.service, then modify the paths and username below

[Unit]
Description=Nono CoWork - Feishu Bot
After=network.target

[Service]
Type=simple
# ↓ Change to your username
User=your-username
# ↓ Change to your project path
WorkingDirectory=/home/your-username/nono-cowork
ExecStart=/home/your-username/nono-cowork/.venv/bin/python -m src.channels.feishu
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
