# feishu-bot.service
# Nono CoWork - Feishu Bot systemd service
#
# Setup:
#   1. Edit this file: fill in your username and project path below
#   2. sudo cp feishu-bot.service /etc/systemd/system/
#   3. sudo systemctl daemon-reload
#   4. sudo systemctl enable --now feishu-bot
#
# Manage:
#   sudo systemctl status feishu-bot       # check status
#   sudo journalctl -u feishu-bot -f       # view logs
#   sudo systemctl restart feishu-bot      # restart after code update

[Unit]
Description=Nono CoWork - Feishu Bot
After=network.target

[Service]
Type=simple
User=your-username
WorkingDirectory=/home/your-username/nono-cowork
ExecStart=/home/your-username/.local/bin/uv run python -m src.channels.feishu
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
Environment=HOME=/home/your-username

[Install]
WantedBy=multi-user.target
